# build/etc/supervisor/conf.d/20-apache.conf
[program:apache]
command=/setup/cryptomedic/start-apache
redirect_stderr=true
