# .devcontainer/build/etc/supervisor/conf.d/mysql.conf
[program:mysqld]
# command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe --pid-file=/var/run/mysqld/mysqld.pid
command=/usr/sbin/mysqld --user=mysql --skip-name-resolve=1 --skip-networking=1 --defaults-extra-file=/etc/mysql/mysql.cnf
user=mysql
redirect_stderr=true

autorestart=false
startretries=0
